# pi-bluetooth/files/pi-bluetooth.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org

#start before bluetoothd (62)
START=61
USE_PROCD=1

start_service() {
	/usr/bin/btuart && /usr/bin/bthelper hci0
}
